# deb/DEBIAN/prerm
#!/bin/sh
systemctl stop kafka.service || echo "Kafka was not running"
systemctl disable kafka.service
